# files/var/lib/pacman/local/eos-bash-shared-24.16.25-1/files
%FILES%
etc/
etc/eos-pacdiff.conf
etc/eos-script-lib-yad.conf
etc/eos-sendlog.conf
usr/
usr/bin/
usr/bin/CalArgConv
usr/bin/ChangeDisplayResolution
usr/bin/RunInTerminal
usr/bin/RunInTerminalEx
usr/bin/RunInTerminalOpt
usr/bin/curl-exit-code-to-string
usr/bin/device-info
usr/bin/eos-FindAppIcon
usr/bin/eos-connection-checker
usr/bin/eos-download-wallpapers
usr/bin/eos-github2gitlab
usr/bin/eos-kbd-set
usr/bin/eos-kernel-nvidia-update-check
usr/bin/eos-pacdiff
usr/bin/eos-reboot-required2
usr/bin/eos-reboot-required3
usr/bin/eos-run-cmd-with-su
usr/bin/eos-select-file-server
usr/bin/eos-sendlog
usr/bin/eos-sendlog-helper
usr/bin/eos-shifttime
usr/bin/eos-sleep-counter
usr/bin/eos-update
usr/bin/eos-update-extras
usr/bin/eos-waiting-indicator
usr/bin/eos-wallpaper-set
usr/bin/grub-colour-chooser
usr/bin/paccache-service-manager
usr/bin/pkg-name-components
usr/bin/show-location-info
usr/bin/su-c_wrapper
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/eos-reboot-required.service
usr/lib/systemd/system/eos-reboot-required.timer
usr/share/
usr/share/applications/
usr/share/applications/eos-update.desktop
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/device-info
usr/share/bash-completion/completions/eos-update
usr/share/bash-completion/completions/show-location-info
usr/share/endeavouros/
usr/share/endeavouros/scripts/
usr/share/endeavouros/scripts/eos-script-lib-yad
usr/share/endeavouros/scripts/ksetwallpaper.py
usr/share/libalpm/
usr/share/libalpm/hooks/
usr/share/libalpm/hooks/eos-reboot-required.hook

%BACKUP%
etc/eos-pacdiff.conf	4d657dbee27bc530b247636072bc7759
etc/eos-script-lib-yad.conf	3c8e646e5593ac1dc1596226f1e51a49
etc/eos-sendlog.conf	5880b45c4bf7f82513dd4c200fffcf1c

