# files/var/lib/pacman/local/libteam-1.32-2/files
%FILES%
usr/
usr/bin/
usr/bin/bond2team
usr/bin/teamd
usr/bin/teamdctl
usr/bin/teamnl
usr/include/
usr/include/team.h
usr/include/teamdctl.h
usr/lib/
usr/lib/libteam.so
usr/lib/libteam.so.5
usr/lib/libteam.so.5.6.2
usr/lib/libteamdctl.so
usr/lib/libteamdctl.so.0
usr/lib/libteamdctl.so.0.1.6
usr/lib/pkgconfig/
usr/lib/pkgconfig/libteam.pc
usr/lib/pkgconfig/libteamdctl.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/teamd@.service
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/teamd.conf
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/bond2team.1.gz
usr/share/man/man5/
usr/share/man/man5/teamd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/teamd.8.gz
usr/share/man/man8/teamdctl.8.gz
usr/share/man/man8/teamnl.8.gz

