# files/var/lib/pacman/local/plocate-1.1.23-1/files
%FILES%
etc/
etc/updatedb.conf
usr/
usr/bin/
usr/bin/locate
usr/bin/mlocate
usr/bin/plocate
usr/bin/plocate-build
usr/bin/updatedb
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/plocate-updatedb.service
usr/lib/systemd/system/plocate-updatedb.timer
usr/lib/systemd/system/timers.target.wants/
usr/lib/systemd/system/timers.target.wants/plocate-updatedb.timer
usr/lib/sysusers.d/
usr/lib/sysusers.d/plocate.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/plocate.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/plocate.1.gz
usr/share/man/man5/
usr/share/man/man5/updatedb.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/plocate-build.8.gz
usr/share/man/man8/updatedb.8.gz
var/
var/lib/
var/lib/plocate/
var/lib/plocate/CACHEDIR.TAG

%BACKUP%
etc/updatedb.conf	e85faa217d111b908b7c1e0e9321b7f1

