# files/var/lib/pacman/local/fwupd-2.0.3-1/files
%FILES%
etc/
etc/fwupd/
etc/fwupd/bios-settings.d/
etc/fwupd/bios-settings.d/README.md
etc/fwupd/fwupd.conf
etc/fwupd/remotes.d/
etc/fwupd/remotes.d/lvfs-testing.conf
etc/fwupd/remotes.d/lvfs.conf
etc/fwupd/remotes.d/vendor-directory.conf
etc/grub.d/
etc/grub.d/35_fwupd
etc/pki/
etc/pki/fwupd-metadata/
etc/pki/fwupd-metadata/GPG-KEY-Linux-Foundation-Metadata
etc/pki/fwupd-metadata/GPG-KEY-Linux-Vendor-Firmware-Service
etc/pki/fwupd-metadata/LVFS-CA.pem
etc/pki/fwupd/
etc/pki/fwupd/GPG-KEY-Linux-Foundation-Firmware
etc/pki/fwupd/GPG-KEY-Linux-Vendor-Firmware-Service
etc/pki/fwupd/LVFS-CA.pem
usr/
usr/bin/
usr/bin/fwupd-dbxtool
usr/bin/fwupdmgr
usr/bin/fwupdtool
usr/include/
usr/include/fwupd-3/
usr/include/fwupd-3/fwupd.h
usr/include/fwupd-3/libfwupd/
usr/include/fwupd-3/libfwupd/fwupd-bios-setting.h
usr/include/fwupd-3/libfwupd/fwupd-build.h
usr/include/fwupd-3/libfwupd/fwupd-client-sync.h
usr/include/fwupd-3/libfwupd/fwupd-client.h
usr/include/fwupd-3/libfwupd/fwupd-codec.h
usr/include/fwupd-3/libfwupd/fwupd-common.h
usr/include/fwupd-3/libfwupd/fwupd-device.h
usr/include/fwupd-3/libfwupd/fwupd-enums.h
usr/include/fwupd-3/libfwupd/fwupd-error.h
usr/include/fwupd-3/libfwupd/fwupd-plugin.h
usr/include/fwupd-3/libfwupd/fwupd-release.h
usr/include/fwupd-3/libfwupd/fwupd-remote.h
usr/include/fwupd-3/libfwupd/fwupd-report.h
usr/include/fwupd-3/libfwupd/fwupd-request.h
usr/include/fwupd-3/libfwupd/fwupd-security-attr.h
usr/include/fwupd-3/libfwupd/fwupd-version.h
usr/lib/
usr/lib/fwupd-2.0.3/
usr/lib/fwupd-2.0.3/libfu_plugin_flashrom.so
usr/lib/fwupd-2.0.3/libfu_plugin_modem_manager.so
usr/lib/fwupd-2.0.3/libfwupdengine.so
usr/lib/fwupd-2.0.3/libfwupdplugin.so
usr/lib/fwupd-2.0.3/libfwupdutil.so
usr/lib/fwupd/
usr/lib/fwupd/fwupd
usr/lib/fwupd/fwupd-detect-cet
usr/lib/girepository-1.0/
usr/lib/girepository-1.0/Fwupd-2.0.typelib
usr/lib/libfwupd.so
usr/lib/libfwupd.so.3
usr/lib/libfwupd.so.3.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/fwupd.pc
usr/lib/systemd/
usr/lib/systemd/system-shutdown/
usr/lib/systemd/system-shutdown/fwupd.shutdown
usr/lib/systemd/system/
usr/lib/systemd/system/fwupd-refresh.service
usr/lib/systemd/system/fwupd-refresh.timer
usr/lib/systemd/system/fwupd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/fwupd.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/fwupdmgr
usr/share/bash-completion/completions/fwupdtool
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/org.freedesktop.fwupd.xml
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.freedesktop.fwupd.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.freedesktop.fwupd.conf
usr/share/doc/
usr/share/doc/fwupd/
usr/share/doc/fwupd/hsi.html
usr/share/doc/fwupd/index.html
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/fwupdmgr.fish
usr/share/fwupd/
usr/share/fwupd/add_capsule_header.py
usr/share/fwupd/firmware_packager.py
usr/share/fwupd/host-emulate.d/
usr/share/fwupd/host-emulate.d/thinkpad-p1-iommu.json.gz
usr/share/fwupd/install_dell_bios_exe.py
usr/share/fwupd/metainfo/
usr/share/fwupd/metainfo/org.freedesktop.fwupd.remotes.lvfs-testing.metainfo.xml
usr/share/fwupd/metainfo/org.freedesktop.fwupd.remotes.lvfs.metainfo.xml
usr/share/fwupd/quirks.d/
usr/share/fwupd/quirks.d/builtin.quirk.gz
usr/share/fwupd/remotes.d/
usr/share/fwupd/remotes.d/fwupd-tests.conf
usr/share/fwupd/remotes.d/vendor/
usr/share/fwupd/remotes.d/vendor/firmware/
usr/share/fwupd/remotes.d/vendor/firmware/README.md
usr/share/fwupd/simple_client.py
usr/share/fwupd/uefi-capsule-ux.tar.xz
usr/share/gir-1.0/
usr/share/gir-1.0/Fwupd-2.0.gir
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/128x128/
usr/share/icons/hicolor/128x128/apps/
usr/share/icons/hicolor/128x128/apps/org.freedesktop.fwupd.png
usr/share/icons/hicolor/64x64/
usr/share/icons/hicolor/64x64/apps/
usr/share/icons/hicolor/64x64/apps/org.freedesktop.fwupd.png
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/org.freedesktop.fwupd.svg
usr/share/locale/
usr/share/locale/af/
usr/share/locale/af/LC_MESSAGES/
usr/share/locale/af/LC_MESSAGES/fwupd.mo
usr/share/locale/ast/
usr/share/locale/ast/LC_MESSAGES/
usr/share/locale/ast/LC_MESSAGES/fwupd.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/fwupd.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/fwupd.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/fwupd.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/fwupd.mo
usr/share/locale/en_GB/
usr/share/locale/en_GB/LC_MESSAGES/
usr/share/locale/en_GB/LC_MESSAGES/fwupd.mo
usr/share/locale/eo/
usr/share/locale/eo/LC_MESSAGES/
usr/share/locale/eo/LC_MESSAGES/fwupd.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/fwupd.mo
usr/share/locale/eu/
usr/share/locale/eu/LC_MESSAGES/
usr/share/locale/eu/LC_MESSAGES/fwupd.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/fwupd.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/fwupd.mo
usr/share/locale/fur/
usr/share/locale/fur/LC_MESSAGES/
usr/share/locale/fur/LC_MESSAGES/fwupd.mo
usr/share/locale/gl/
usr/share/locale/gl/LC_MESSAGES/
usr/share/locale/gl/LC_MESSAGES/fwupd.mo
usr/share/locale/he/
usr/share/locale/he/LC_MESSAGES/
usr/share/locale/he/LC_MESSAGES/fwupd.mo
usr/share/locale/hi/
usr/share/locale/hi/LC_MESSAGES/
usr/share/locale/hi/LC_MESSAGES/fwupd.mo
usr/share/locale/hr/
usr/share/locale/hr/LC_MESSAGES/
usr/share/locale/hr/LC_MESSAGES/fwupd.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/fwupd.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/fwupd.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/fwupd.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/fwupd.mo
usr/share/locale/ka/
usr/share/locale/ka/LC_MESSAGES/
usr/share/locale/ka/LC_MESSAGES/fwupd.mo
usr/share/locale/kk/
usr/share/locale/kk/LC_MESSAGES/
usr/share/locale/kk/LC_MESSAGES/fwupd.mo
usr/share/locale/ko/
usr/share/locale/ko/LC_MESSAGES/
usr/share/locale/ko/LC_MESSAGES/fwupd.mo
usr/share/locale/ky/
usr/share/locale/ky/LC_MESSAGES/
usr/share/locale/ky/LC_MESSAGES/fwupd.mo
usr/share/locale/lt/
usr/share/locale/lt/LC_MESSAGES/
usr/share/locale/lt/LC_MESSAGES/fwupd.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/fwupd.mo
usr/share/locale/oc/
usr/share/locale/oc/LC_MESSAGES/
usr/share/locale/oc/LC_MESSAGES/fwupd.mo
usr/share/locale/pa/
usr/share/locale/pa/LC_MESSAGES/
usr/share/locale/pa/LC_MESSAGES/fwupd.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/fwupd.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/fwupd.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/fwupd.mo
usr/share/locale/ro/
usr/share/locale/ro/LC_MESSAGES/
usr/share/locale/ro/LC_MESSAGES/fwupd.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/fwupd.mo
usr/share/locale/si/
usr/share/locale/si/LC_MESSAGES/
usr/share/locale/si/LC_MESSAGES/fwupd.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/fwupd.mo
usr/share/locale/sl/
usr/share/locale/sl/LC_MESSAGES/
usr/share/locale/sl/LC_MESSAGES/fwupd.mo
usr/share/locale/sr/
usr/share/locale/sr/LC_MESSAGES/
usr/share/locale/sr/LC_MESSAGES/fwupd.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/fwupd.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/fwupd.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/fwupd.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/fwupd.mo
usr/share/locale/zh_TW/
usr/share/locale/zh_TW/LC_MESSAGES/
usr/share/locale/zh_TW/LC_MESSAGES/fwupd.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/fwupd-dbxtool.1.gz
usr/share/man/man1/fwupdmgr.1.gz
usr/share/man/man1/fwupdtool.1.gz
usr/share/man/man5/
usr/share/man/man5/fwupd-remotes.d.5.gz
usr/share/man/man5/fwupd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/fwupd-refresh.service.8.gz
usr/share/metainfo/
usr/share/metainfo/org.freedesktop.fwupd.metainfo.xml
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/org.freedesktop.fwupd.policy
usr/share/polkit-1/rules.d/
usr/share/polkit-1/rules.d/org.freedesktop.fwupd.rules
usr/share/vala/
usr/share/vala/vapi/
usr/share/vala/vapi/fwupd.deps
usr/share/vala/vapi/fwupd.vapi

%BACKUP%
etc/fwupd/fwupd.conf	37a4f57716b9e3ebd0716a53306665b7
etc/fwupd/remotes.d/lvfs-testing.conf	ff89b44a3deff82b028461f09dd056df
etc/fwupd/remotes.d/lvfs.conf	2428a78af3b78452023f4c9a32855037
etc/fwupd/remotes.d/vendor-directory.conf	444a68b9510616e1fe30c6c7435f9cd9

