# files/var/lib/pacman/local/libcanberra-1:0.30+r2+gc0620e4-4/files
%FILES%
etc/
etc/X11/
etc/X11/xinit/
etc/X11/xinit/xinitrc.d/
etc/X11/xinit/xinitrc.d/40-libcanberra-gtk-module.sh
usr/
usr/bin/
usr/bin/canberra-boot
usr/bin/canberra-gtk-play
usr/include/
usr/include/canberra-gtk.h
usr/include/canberra.h
usr/lib/
usr/lib/gnome-settings-daemon-3.0/
usr/lib/gnome-settings-daemon-3.0/gtk-modules/
usr/lib/gnome-settings-daemon-3.0/gtk-modules/canberra-gtk-module.desktop
usr/lib/gtk-2.0/
usr/lib/gtk-2.0/modules/
usr/lib/gtk-2.0/modules/libcanberra-gtk-module.so
usr/lib/gtk-3.0/
usr/lib/gtk-3.0/modules/
usr/lib/gtk-3.0/modules/libcanberra-gtk-module.so
usr/lib/gtk-3.0/modules/libcanberra-gtk3-module.so
usr/lib/libcanberra-0.30/
usr/lib/libcanberra-0.30/libcanberra-alsa.la
usr/lib/libcanberra-0.30/libcanberra-alsa.so
usr/lib/libcanberra-0.30/libcanberra-multi.la
usr/lib/libcanberra-0.30/libcanberra-multi.so
usr/lib/libcanberra-0.30/libcanberra-null.la
usr/lib/libcanberra-0.30/libcanberra-null.so
usr/lib/libcanberra-0.30/libcanberra-pulse.la
usr/lib/libcanberra-0.30/libcanberra-pulse.so
usr/lib/libcanberra-gtk.so
usr/lib/libcanberra-gtk.so.0
usr/lib/libcanberra-gtk.so.0.1.9
usr/lib/libcanberra-gtk3.so
usr/lib/libcanberra-gtk3.so.0
usr/lib/libcanberra-gtk3.so.0.1.9
usr/lib/libcanberra.so
usr/lib/libcanberra.so.0
usr/lib/libcanberra.so.0.2.5
usr/lib/pkgconfig/
usr/lib/pkgconfig/libcanberra-gtk.pc
usr/lib/pkgconfig/libcanberra-gtk3.pc
usr/lib/pkgconfig/libcanberra.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/canberra-system-bootup.service
usr/lib/systemd/system/canberra-system-shutdown-reboot.service
usr/lib/systemd/system/canberra-system-shutdown.service
usr/share/
usr/share/gtk-doc/
usr/share/gtk-doc/html/
usr/share/gtk-doc/html/libcanberra/
usr/share/gtk-doc/html/libcanberra/ch01.html
usr/share/gtk-doc/html/libcanberra/home.png
usr/share/gtk-doc/html/libcanberra/index.html
usr/share/gtk-doc/html/libcanberra/left-insensitive.png
usr/share/gtk-doc/html/libcanberra/left.png
usr/share/gtk-doc/html/libcanberra/libcanberra-canberra-gtk.html
usr/share/gtk-doc/html/libcanberra/libcanberra-canberra.html
usr/share/gtk-doc/html/libcanberra/libcanberra.devhelp2
usr/share/gtk-doc/html/libcanberra/right-insensitive.png
usr/share/gtk-doc/html/libcanberra/right.png
usr/share/gtk-doc/html/libcanberra/style.css
usr/share/gtk-doc/html/libcanberra/up-insensitive.png
usr/share/gtk-doc/html/libcanberra/up.png
usr/share/vala/
usr/share/vala/vapi/
usr/share/vala/vapi/libcanberra-gtk.vapi
usr/share/vala/vapi/libcanberra.vapi

