# files/var/lib/pacman/local/cups-browsed-2.1.0-1/files
%FILES%
etc/
etc/cups/
etc/cups/cups-browsed.conf
usr/
usr/bin/
usr/bin/cups-browsed
usr/lib/
usr/lib/cups/
usr/lib/cups/backend/
usr/lib/cups/backend/implicitclass
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/cups-browsed.service
usr/share/
usr/share/doc/
usr/share/doc/cups-browsed/
usr/share/doc/cups-browsed/ABOUT-NLS
usr/share/doc/cups-browsed/AUTHORS
usr/share/doc/cups-browsed/CHANGES-1.x.md
usr/share/doc/cups-browsed/CHANGES.md
usr/share/doc/cups-browsed/CONTRIBUTING.md
usr/share/doc/cups-browsed/COPYING
usr/share/doc/cups-browsed/DEVELOPING.md
usr/share/doc/cups-browsed/INSTALL
usr/share/doc/cups-browsed/LICENSE
usr/share/doc/cups-browsed/NOTICE
usr/share/doc/cups-browsed/README.md
usr/share/licenses/
usr/share/licenses/cups-browsed/
usr/share/licenses/cups-browsed/COPYING
usr/share/licenses/cups-browsed/NOTICE
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/cups-browsed.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/cups-browsed.8.gz

%BACKUP%
etc/cups/cups-browsed.conf	3f33675be6f308c1c52839186919b115

